# kbdremap.service
[Unit]
After=network.service

[Service]
# User=
ExecStart=/usr/local/bin/kbdremap.sh

[Install]
WantedBy=default.target
